# docs/deploy/fluidcal-db-backup.timer
[Unit]
Description=Run the fluidcal database backup nightly

[Timer]
OnCalendar=*-*-* 03:30:00
Persistent=true
RandomizedDelaySec=10m

[Install]
WantedBy=timers.target
